# .config/systemd/user/psd-resync.timer.d/frequency.conf
[Unit]
Description=Custom timer for Profile Sync Daemon

[Timer]
OnUnitActiveSec=
OnUnitActiveSec=10min
